# deploy/kalshi-bot.service
[Unit]
Description=Kalshi Weather Arbitrage Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/weather-bot
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 main.py
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=kalshi-bot

[Install]
WantedBy=multi-user.target
